# TV/webserver.service
[Unit]
Description=Running a webserver on locahost
After=network.target

[Service]
ExecStart=/home/pi/server/env/bin/python3 -u webserver.py
WorkingDirectory=/home/pi/server
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target